# docker-images/gluster-nfs-s6/rootfs/etc/services.d/gluster/finish
#!/usr/bin/execlineb -S0

# Bring the container down if the process exits
s6-svscanctl -t /var/run/s6/services
